# modules/unicorn/templates/headers/_redhat.sh.erb
#!/bin/sh
# chkconfig: 2345 95 20
# description: Unicorn script
# Manage the unicorn service as a non privileged user (at least the user not need admin)
# processname: unicorn
