# modules/unicorn/templates/unicorn.service.erb
[Unit]
Description=Unicorn Server

[Service]
User=<%= scope['user::deployer::username'] %>
Group=<%= scope['user::deployer::group'] %>
WorkingDirectory=<%= scope['vars::apps_root'] %>/<%= @fqdn %>/current
Environment=RAILS_ENV=<%= @virtual == 'virtualbox' ? 'sandbox' : 'production' %>
SyslogIdentifier=unicorn
PIDFile=/run/unicorn.pid

ExecStart=<%= scope['user::deployer::home'] %>/.rbenv/shims/bundle exec \
  <%= scope['vars::apps_root'] %>/<%= @fqdn %>/current/bin/unicorn -D -c \
  <%= scope['vars::apps_root'] %>/<%= @fqdn %>/current/config/unicorn.rb
ExecStop=/usr/bin/kill -QUIT $MAINPID
ExecReload=/usr/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
