# modules/meteor/templates/headers/_redhat.sh.erb
#!/bin/sh
# chkconfig: 2345 99 20
# description: Meteor script
# Manage the meteor service
# processname: meteor
