# modules/sidekiq/templates/sidekiq.service.erb
[Unit]
Description=Sidekiq
Requires=redis.service
After=redis.service syslog.target network.target

[Service]
Type=simple
User=<%= scope['user::deployer::username'] %>
Group=<%= scope['user::deployer::group'] %>
UMask=0002
WorkingDirectory=<%= scope['vars::apps_root'] %>/<%= @fqdn %>/current

ExecStart=<%= scope['user::deployer::home'] %>/.rbenv/shims/bundle exec \
  "sidekiq -e <%= @productname == 'VirtualBox' ? 'sandbox' : 'production' %> \
  -C config/sidekiq.yml -L log/sidekiq.log "

# if we crash, restart
RestartSec=1
Restart=on-failure

[Install]
WantedBy=multi-user.target
