# modules/oracle/templates/headers/_debian.sh.erb
### BEGIN INIT INFO
# Provides:          oracle
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: oracle
# Description:       oracle server autostart-script
### END INIT INFO
