# modules/oracle/templates/headers/_redhat.sh.erb
#!/bin/sh
# chkconfig: 2345 94 20
# description: Oracle script
# Manage the oracle service
# processname: oracle
