# modules/sidekiq/templates/headers/_debian.sh.erb
### BEGIN INIT INFO
# Provides:          sidekiq
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: sidekiq
# Description:       sidekiq server autostart-script
### END INIT INFO
