# modules/sidekiq/templates/headers/_redhat.sh.erb
#!/bin/sh
# chkconfig: 2345 95 20
# description: Sidekiq script
# Manage the sidekiq service as a non privileged user (at least the user not need admin)
# processname: sidekiq
