# modules/sidekiq/templates/headers/_redhat.sh.erb
#!/bin/sh
# chkconfig: 2345 99 20
# description: Sidekiq script
# Manage the sidekiq service
# processname: sidekiq
